# extras/systemd/he-update.timer
[Unit]
Description=Run he-update every 5 minutes

[Timer]
OnBootSec=50
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
